# udev-rules/81.power.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/share/slimbookamdcontroller/src/applyconfig.py"
